# files/systemd/automatic-update.service
[Unit]
Description=Automatic weekly upgrade

[Service]
ExecStart=/root/system-scripts/automatic-update.sh