# systemd/firstboot/install-xanmod-kernel.service
[Unit]
Description=Install Xanmod Kernel on First Boot
ConditionFirstBoot=yes

[Service]
Type=oneshot
ExecStart=/root/first_boot.sh

[Install]
WantedBy=multi-user.target